# ansible/files/etc/systemd/system/physlock.service.j2
# {{ ansible_managed }}

[Unit]
Description=Lock session using physlock

[Service]
User={{ user }}
ExecStart=/usr/bin/physlock
[Install]
WantedBy=sleep.target  
