# etc/systemd/system/physlock@.service
[Unit]
Description=Lock session using physlock

[Service]
User=%i
ExecStart=/usr/bin/physlock
[Install]
WantedBy=sleep.target  